# systemd/tgbot2-worker.service
[Unit]
Description=Sonvia IVR Bot - Campaign Worker
After=network.target postgresql.service tgbot4-webhook.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tgbot4/dialer
ExecStart=/usr/bin/python3 campaign_worker.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
