# systemd/tgbot2-webhook.service
[Unit]
Description=Sonvia IVR Bot - Webhook Server (Port 8003)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tgbot4/dialer
ExecStart=/usr/bin/python3 webhook_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
